# systemd/agent-usage-poll.service
[Unit]
Description=Refresh agent usage state
After=default.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/agent-usage-widget/.env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/bin/agent-usage-poll
